# systemd/saferoads.timer
[Unit]
Description=Run SafeRoads container hourly (Europe/London)

[Timer]
OnCalendar=hourly
Timezone=Europe/London
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
